# systemd/tg-digest-bot.service
[Unit]
Description=Telegram Digest Bot
After=network.target

[Service]
Type=simple
User=tg-digest-bot
Group=tg-digest-bot
WorkingDirectory=/home/tg-digest-bot/tg-digest-bot
Environment=PYTHONPATH=/home/tg-digest-bot/tg-digest-bot/src
ExecStart=/home/tg-digest-bot/tg-digest-bot/.venv/bin/python -m tg_digest_bot
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
